# scripts/bryn/cron/jflory.cron
@reboot /home/jflory/git/infrastructure/scripts/bryn/restore_tmux_sessions_boot.sh
