# scripts/archived/fiction/cron/root.cron
00 05 15 * * /home/jflory/git/infrastructure/scripts/fiction/renew_ssl_certs.sh
